# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:backend]
directory=/workspace/backend
command=uvicorn app.main:app --host 127.0.0.1 --port 8000
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=CROPTWIN_DISEASE_ARTIFACT_DIR="/workspace/backend/model_artifacts/croptwin_disease"

[program:frontend]
directory=/workspace
command=streamlit run frontend/app.py --server.address 0.0.0.0 --server.port 7860 --server.headless true
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=CROPTWIN_API_BASE_URL="http://127.0.0.1:8000"